# data/bec5.service.in
[Unit]
Description=BEC 5 web interface
After=syslog.target
Require=bec5-dbus.service

[Service]
ExecStart=/usr/bin/uwsgi --ini @UWSGI_INI_PATH@/bec5.ini
Restart=always
Type=notify
StandardError=syslog
NotifyAccess=main

[Install]
WantedBy=multi-user.target
